# monitor.service
[Unit]
Description=Monitor bot directory

[Service]
ExecStart=/home/pi/bot/monitor.sh

[Install]
WantedBy=multi-user.target
